# qbit.timer
[Unit]
Description=Run qbittorent after 1min of boot
[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target